# launchd/com.reviewinator.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.reviewinator</string>
	<key>ProgramArguments</key>
	<array>
		<string>/Users/mmclane/.local/bin/uv</string>
		<string>run</string>
		<string>python</string>
		<string>-m</string>
		<string>reviewinator</string>
	</array>
	<key>WorkingDirectory</key>
	<string>/Users/mmclane/repos/personal/reviewinator</string>
	<key>RunAtLoad</key>
	<true/>
	<key>KeepAlive</key>
	<true/>
	<key>StandardOutPath</key>
	<string>/Users/mmclane/Library/Logs/reviewinator/stdout.log</string>
	<key>StandardErrorPath</key>
	<string>/Users/mmclane/Library/Logs/reviewinator/stderr.log</string>
</dict>
</plist>
